# vcenter/crontab
*/5 * * * * python3 /app/generate.py > $OUTPUT_PATH