# vcenter/crontab
*/5 * * * * mkdir -p $OUTPUT_PATH python generate.py > $OUTPUT_PATH