# vcenter/crontab
*/5 * * * * python3 generate.py > $OUTPUT_PATH