# vcenter/crontab
*/5 * * * * python generate.py > $OUTPUT_PATH